# conf/etc/cron.hourly/ntp-update
#!/bin/sh

/usr/sbin/ntpdate pool.ntp.org 2>&1 >/dev/null
